# packer/files/etc/systemd/system/haproxy.service.d/extra.conf
[Service]
EnvironmentFile=/etc/sysconfig/haproxy
Restart=on-failure
RestartSec=5